# systemd/careersignal.service
[Unit]
Description=CareerSignal Web App (FastAPI)
After=network.target

[Service]
Type=simple
User=careersignal
Group=careersignal
WorkingDirectory=/home/careersignal/careersignal
Environment=PYTHONUNBUFFERED=1
# Load environment (secrets/SMTP/etc.)
EnvironmentFile=-/etc/careersignal.env

ExecStart=/home/careersignal/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8010
Restart=always
RestartSec=2

# Hardening (safe defaults)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/home/careersignal/careersignal/data

[Install]
WantedBy=multi-user.target
